# systemd/better-rm-trash-cleanup.service
[Unit]
Description=Clean old files from better-rm trash directories
Documentation=man:better-rm(1)
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/trash-cleanup.sh
StandardOutput=journal
StandardError=journal
Nice=19
IOSchedulingClass=idle

# Security hardening
PrivateTmp=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home /root

[Install]
WantedBy=multi-user.target
