# systemd/better-rm-trash-cleanup.service
[Unit]
Description=Run better_rm trash cleanup weekly
Documentation=man:better_rm(1)
Requires=better_rm-trash-cleanup.service

[Timer]
# Run every Sunday at 3:00 AM
OnCalendar=weekly
# Also run 15 minutes after boot if a scheduled run was missed
OnBootSec=15min
# Randomize start time by up to 1 hour to prevent load spikes
RandomizedDelaySec=1h
# Save last trigger time
Persistent=true

[Install]
WantedBy=timers.target
